# config/com.tartarus.linear-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tartarus.linear-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/curl</string>
        <string>-s</string>
        <string>-X</string>
        <string>POST</string>
        <string>-H</string>
        <string>Content-Type: application/json</string>
        <string>-d</string>
        <string>{"includeCompleted": true}</string>
        <string>http://localhost:3005/api/integrations/linear/sync</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/tartarus-linear-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tartarus-linear-sync.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
